# etc/services.d/mtproxy/run
#!/usr/bin/with-contenv sh

#./usr/bin/start_mtproxy.sh
bash /usr/bin/start_mtproxy.sh